# targets/bandit/src/cronjob_bandit22.cron
* * * * * bandit22 /usr/bin/cronjob_bandit22.sh
