# targets/bandit/src/cronjob_bandit24.cron
* * * * * root /usr/bin/cronjob_bandit24.sh
