# targets/bandit/src/cronjob_bandit23.cron
* * * * * bandit23 /usr/bin/cronjob_bandit23.sh
